# gateway.service
[Unit]
Decription=gateway

[Service]
ExecStart=/root/start_gateway.sh

[Install]
WantedBy=multi-user.target

